# nginx/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /var/opt/webops/nginx/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
